# ansible/roles/flame/templates/flame.service.j2
[Unit]
Description=Service for flame

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ flame_app_directory }}
ExecStart={{ docker_compose_binary }} up -d
ExecStop={{ docker_compose_binary }} stop
ExecReload={{ docker_compose_binary }} restart

[Install]
WantedBy=multi-user.target
