# ansible/roles/transmission/templates/transmission.service.j2
[Unit]
Description=Service for transmission

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ transmission_app_directory }}
ExecStart={{ docker_compose_binary }} up -d
ExecStop={{ docker_compose_binary }} stop
ExecReload={{ docker_compose_binary }} restart

[Install]
WantedBy=multi-user.target
