# ansible/roles/yamtrack/templates/yamtrack.service.j2
[Unit]
Description=Service for yamtrack
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ yamtrack_app_directory }}
ExecStart={{ docker_compose_binary }} up -d
ExecStop={{ docker_compose_binary }} stop
ExecReload={{ docker_compose_binary }} restart

[Install]
WantedBy=multi-user.target
